# systemd/kdrive-rclone.service
[Unit]
Description=Mount kDrive with rclone
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/rclone mount kdrive: /mnt/kdrive \
  --allow-other \
  --vfs-cache-mode full \
  --dir-cache-time 1h
ExecStop=/bin/fusermount -u /mnt/kdrive
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
